# app/curator/crontab
@daily /curator/run-curator.sh >> /curator/log.log 2>&1
#
